# image/runit/ui/run
#!/usr/bin/env bash

cd /ui
exec npx http-server -p ${PORT:-3000} build