# image/runit/ui/run
#!/usr/bin/env bash

cd /ui
exec npx serve -l ${PORT:-3000} --config config/serve.json